# internal/install/systemd/marc-process.service.tmpl
[Unit]
Description=marc-server process daemon (MinIO poll, denoise, ClickHouse ingest)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{.BinaryPath}} process --config /etc/marc/server.toml
Restart=on-failure
RestartSec=5s
DynamicUser=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/marc /var/log
NoNewPrivileges=true
PrivateTmp=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
