# internal/install/systemd/marc-generate.timer.tmpl
[Unit]
Description=Hourly timer for marc-server generate

[Timer]
OnCalendar=hourly
Persistent=true
Unit=marc-generate.service

[Install]
WantedBy=timers.target
