# internal/install/templates/marc-proxy.service.tmpl
[Unit]
Description=marc proxy daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{.BinaryPath}} proxy --config {{.ConfigPath}}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
