# internal/install/templates/marc-ship.service.tmpl
[Unit]
Description=marc ship daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{.BinaryPath}} ship --config {{.ConfigPath}}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
{{if .User}}User={{.User}}
{{end}}{{if .Group}}Group={{.Group}}
{{end}}NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
{{if .HomeDir}}ReadWritePaths={{.HomeDir}}/.marc
{{end}}

[Install]
WantedBy=default.target
